# .config/systemd/user/neuter-kube-context.timer
[Unit]
Description=Neuters kube context

[Timer]
Unit=neuter-kube-context.service
OnCalendar=*:0/15

[Install]
WantedBy=timers.target